# deployment/roles/wheelodex/templates/mail-systemd-failure@.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/mail-systemd-failure %i
User=nobody
Group=systemd-journal
